# packaging/aur/pkg/niri-display-manager/usr/lib/systemd/user/niri-display-manager-daemon.service
[Unit]
Description=Niri Display Manager Hotplug Daemon
Documentation=https://github.com/rickycbanks/niri_display_manager
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/niri-display-manager --daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
